# back/localdeploy/cron/crontab
# /etc/cron.d/sco_auth
*/15 * * * * /usr/local/bin/cron_script >> /var/log/crontab.log 2>&1
